# .docker/supervisor/conf/laravel-worker.conf
[supervisord]
nodaemon=true
logfile=/dev/null

[program:laravel-worker]

   process_name=%(program_name)s_%(process_num)02d

   command=php artisan queue:work -v --sleep=3 --tries=3

   autostart=true

   autorestart=true

   user=www-data

   numprocs=1

   redirect_stderr=true

   stdout_logfile=storage/logs/worker.log